# dotfiles/autostart/.config/autostart/alacritty-runner.desktop
[Desktop Entry]
Exec=alacritty --class Alacritty-Runner
Name=alacritty
Type=Application
