# dotfiles/autostart/.config/autostart/alacritty-notes.desktop
[Desktop Entry]
Exec=alacritty --class Alacritty-DigitalGarden --hold -e zsh -l -i -c 'dev digitalgarden'
Name=alacritty
Type=Application
