# dotfiles/autostart/.config/autostart/alacritty.desktop
[Desktop Entry]
Exec=alacritty
Name=alacritty
Type=Application
