# dotfiles/autostart/.config/autostart/alacritty.desktop
[Desktop Entry]
Exec=alacritty --hold -e zsh -l -i -c 'dev open juanibiapina/dotfiles'
Name=alacritty
Type=Application
