# dotfiles/autostart/.config/autostart/alacritty-runner.desktop
[Desktop Entry]
Exec=alacritty --class Alacritty-Runner --hold -e zsh -l -i -c 'dev tmux runner'
Name=alacritty
Type=Application
