# dotfiles/autostart/.config/autostart/dropbox.desktop
[Desktop Entry]
Exec=dropbox
Name=dropbox
Type=Application
