# dotfiles/autostart/.config/autostart/firefox.desktop
[Desktop Entry]
Exec=firefox-developer-edition
Name=Firefox
Type=Application
