# dotfiles/autostart/.config/autostart/firefox-devedition.desktop
[Desktop Entry]
Exec=firefox-devedition
Name=Firefox
Type=Application
