# dotfiles/autostart/.config/autostart/hamsket.desktop
[Desktop Entry]
Exec=hamsket
Name=hamsket
Type=Application
