# dotfiles/autostart/.config/autostart/keepassxc.desktop
[Desktop Entry]
Exec=keepassxc
Name=keepassxc
Type=Application
