# dotfiles/autostart/.config/autostart/nautilus.desktop
[Desktop Entry]
Name=nautilus
Exec=nautilus
Type=Application
