# dotfiles/autostart/.config/autostart/obsidian.desktop
[Desktop Entry]
Exec=obsidian
Name=Obsidian
Type=Application
