# dotfiles/autostart/.config/autostart/pasystray.desktop
[Desktop Entry]
Exec=pasystray
Name=pasystray
Type=Application
