# dotfiles/autostart/.config/autostart/pcloud.desktop
[Desktop Entry]
Exec=pcloud
Name=pcloud
Type=Application
