# dotfiles/autostart/.config/autostart/slack.desktop
[Desktop Entry]
Exec=slack
Name=slack
Type=Application
