# dotfiles/autostart/.config/autostart/spotify.desktop
[Desktop Entry]
Exec=kitty --class Kitty-Spotify --hold zsh -l -i -c 'spotify_player'
Name=spotify
Type=Application
