# dotfiles/autostart/.config/autostart/spotify.desktop
[Desktop Entry]
Exec=spotify
Name=spotify
Type=Application
