# dotfiles/autostart/.config/autostart/udiskie.desktop
[Desktop Entry]
Exec=udiskie --tray
Name=udiskie
Type=Application
