# dotfiles/autostart/.config/autostart/ulauncher.desktop
[Desktop Entry]
Exec=ulauncher --hide-window --no-window-shadow
Name=ulauncher
Type=Application
