# dotfiles/autostart/.config/autostart/vivaldi.desktop
[Desktop Entry]
Exec=vivaldi
Name=Vivaldi
Type=Application
