# dotfiles/autostart/.config/autostart/xbindkeys.desktop
[Desktop Entry]
Name=xbindkeys
Exec=xbindkeys
Type=Application
